# my_todo_app_addon/root/etc/services.d/my_todo_app/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting Node.js server for My To-Do App..."

declare -i port_internal
port_internal=$(bashio::config 'port_internal')

export PORT="${port_internal}"
export DATABASE_PATH="/data/todo.db"

mkdir -p "$(dirname "${DATABASE_PATH}")"

bashio::log.info "Database path: ${DATABASE_PATH}"
bashio::log.info "Server will listen on port: ${PORT}"

NODE_BIN=$(which node)
if [ -z "$NODE_BIN" ]; then
    bashio::log.error "Node.js executable not found in PATH!"
    exit 1
fi

# Stay in the foreground!
exec "$NODE_BIN" /app/backend/server.js
