# my_todo_app_addon/root/etc/services.d/my_todo_app/run
#!/command/with-contenv bash
# DO NOT use /usr/bin/with-contenv — this causes suexec error

set -e

# Logging
echo "[INFO] Starting Node.js server for My To-Do App..."

# Read config
PORT=$(bashio::config 'port_internal')
export PORT="${PORT}"
export DATABASE_PATH="/data/todo.db"

echo "[INFO] Database path: $DATABASE_PATH"
echo "[INFO] Server will listen on port: $PORT"

# Ensure DB directory exists
mkdir -p "$(dirname "$DATABASE_PATH")"

# Start the server as PID 1
exec node /app/backend/server.js
